# com.aaron.remarkable.sync.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.aaron.remarkable.sync</string>
    
    <key>ProgramArguments</key>
    <array>
        <string>/bin/bash</string>
        <string>/Users/aaronmacbaron/Documents/remarkable-mcp/run_remarkable_to_apple_notes.sh</string>
    </array>

    <!-- Run from Home directory to avoid "getcwd" permission errors in Documents -->
    <key>WorkingDirectory</key>
    <string>/Users/aaronmacbaron</string>

    <key>StandardOutPath</key>
    <string>/Users/aaronmacbaron/Documents/remarkable-mcp/logs/launchd.log</string>
    <key>StandardErrorPath</key>
    <string>/Users/aaronmacbaron/Documents/remarkable-mcp/logs/launchd.err</string>

    <!-- Run every 3600 seconds (1 hour) -->
    <key>StartInterval</key>
    <integer>3600</integer>

    <!-- Also run when the computer wakes up/logs in -->
    <key>RunAtLoad</key>
    <true/>
</dict>
</plist>
