# packaging/scripts/postinstall
#!/bin/bash

APP_DIR="/Applications/Living Ink"
APP_PLIST="com.aaron.livingink.plist"

# Disable old launch agent if present (handle both old and new names during transition)
if launchctl print gui/"501"/com.aaron.remarkable.sync >/dev/null 2>&1; then
    launchctl bootout gui/"501" /Library/LaunchAgents/com.aaron.remarkable.sync.plist 2>/dev/null || true
fi
if launchctl print gui/"501"/com.aaron.livingink >/dev/null 2>&1; then
    launchctl bootout gui/"501" /Library/LaunchAgents/"" 2>/dev/null || true
fi
rm -f /Library/LaunchAgents/com.aaron.remarkable.sync.plist
rm -f /Library/LaunchAgents/""

# Copy new plist
cp "$APP_DIR/$APP_PLIST" /Library/LaunchAgents/

# Fix ownership
chown root:wheel /Library/LaunchAgents/""
chmod 644 /Library/LaunchAgents/""

# Load the service
if [ -n "$UID" ]; then
    launchctl bootstrap gui/"501" /Library/LaunchAgents/""
fi

# Get the logged-in user (since installer runs as root)
LOGIN_USER=$(scutil <<< "show State:/Users/ConsoleUser" | awk '/Name :/ && ! /loginwindow/ { print $3 }')

# Fix permissions so the user can run it and write logs
chmod -R 777 "$APP_DIR"

# Wait a moment for the filesystem to settle
sleep 1

# Open as the logged-in user
if [ -n "$LOGIN_USER" ]; then
    # Create Desktop Shortcut
    USER_DESKTOP="/Users/$LOGIN_USER/Desktop"
    if [ -d "$USER_DESKTOP" ]; then
         TARGET_LINK="$USER_DESKTOP/Living Ink Sync.command"
         ln -sf "$APP_DIR/Run Manual Sync.command" "$TARGET_LINK"
         chown "$LOGIN_USER" "$TARGET_LINK"
    fi

    sudo -u "$LOGIN_USER" open "$APP_DIR/docs/SETUP_GUIDE.md"
    sudo -u "$LOGIN_USER" open "$APP_DIR"
fi

exit 0
