# packaging/com.aaron.livingink.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.aaron.livingink</string>
    
    <key>ProgramArguments</key>
    <array>
        <string>/Applications/Living Ink/living-ink</string>
    </array>

    <!-- Run from the App folder so it finds secrets/ and logs/ -->
    <key>WorkingDirectory</key>
    <string>/Applications/Living Ink</string>

    <key>StandardOutPath</key>
    <string>/Applications/Living Ink/logs/launchd.log</string>
    <key>StandardErrorPath</key>
    <string>/Applications/Living Ink/logs/launchd.err</string>

    <!-- Run every 3600 seconds (1 hour) -->
    <key>StartInterval</key>
    <integer>3600</integer>

    <!-- Also run when the computer wakes up/logs in -->
    <key>RunAtLoad</key>
    <true/>
</dict>
